# templates/systemd-user/context-healthcheck.service
[Unit]
Description=Unified Context Healthcheck (user)

[Service]
Type=oneshot
ExecStart=/bin/bash %h/.local/share/context-mesh-foundry/scripts/context_healthcheck.sh --quiet
