# templates/launchd/com.openviking.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.openviking.daemon</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>VIKING_ENABLE_SHELL_MONITOR</key><string>0</string>
    <key>VIKING_ENABLE_OPENCODE_MONITOR</key><string>0</string>
    <key>VIKING_ENABLE_KILO_MONITOR</key><string>0</string>
    <key>VIKING_POLL_INTERVAL_SEC</key><string>180</string>
    <key>VIKING_FAST_POLL_INTERVAL_SEC</key><string>20</string>
    <key>VIKING_IDLE_SLEEP_CAP_SEC</key><string>600</string>
    <key>VIKING_CODEX_SESSION_SCAN_INTERVAL_SEC</key><string>300</string>
    <key>VIKING_CLAUDE_TRANSCRIPT_SCAN_INTERVAL_SEC</key><string>300</string>
    <key>VIKING_ANTIGRAVITY_SCAN_INTERVAL_SEC</key><string>300</string>
    <key>VIKING_SUSPEND_ANTIGRAVITY_WHEN_BUSY</key><string>1</string>
    <key>VIKING_ANTIGRAVITY_BUSY_LS_THRESHOLD</key><string>2</string>
    <key>VIKING_ANTIGRAVITY_INGEST_MODE</key><string>final_only</string>
    <key>VIKING_ANTIGRAVITY_QUIET_SEC</key><string>240</string>
    <key>VIKING_ANTIGRAVITY_MIN_DOC_BYTES</key><string>500</string>
    <key>VIKING_CYCLE_BUDGET_SEC</key><string>8</string>
    <key>VIKING_INDEX_SYNC_MIN_INTERVAL_SEC</key><string>20</string>
    <key>VIKING_ERROR_BACKOFF_MAX_SEC</key><string>30</string>
    <key>VIKING_LOOP_JITTER_SEC</key><string>0.7</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>__SCRIPTS_DIR__/viking_daemon.py</string>
  </array>
  <key>WorkingDirectory</key><string>__SCRIPTS_DIR__</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>__LOG_DIR__/viking_daemon_launchd.log</string>
  <key>StandardErrorPath</key><string>__LOG_DIR__/viking_daemon_launchd.err</string>
</dict>
</plist>
